# ponyboot/var/ftpd/install-k8s.service
[Unit]
Description=Install Kubernetes

[Service]
Type=oneshot
User=root
ExecStart=/install-k8s.sh

[Install]
WantedBy=multi-user.target
